# docker/supervisord.conf
# ============================================================================
#  supervisord.conf — держит дашборд + 8 приложений в одном контейнере.
#  Логи каждого процесса → stdout/stderr контейнера (префикс [program:*]).
#      docker logs -f panel
# ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# --- Дашборд (плитки) :333 — Caddy, статика ---
[program:dashboard]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- cf :3331 — PHP 8.1 + Apache (Cloudflare Panel) ---
[program:cf]
command=apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- cf: фоновый обработчик очереди (был в docker-entrypoint.sh, порт 1000 → 3331) ---
[program:cf-queue]
command=/bin/bash -c "sleep 8; while true; do curl -s -A QueueProcessor 'http://127.0.0.1:3331/queue_processor.php?auth_token=cloudflare_queue_processor_2024' >/dev/null 2>&1 || true; sleep 20; done"
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- cf: фоновый мониторинг доменов (NS/IP/SSL/WHOIS + Telegram) ---
[program:cf-monitor]
command=/bin/bash -c "sleep 25; while true; do curl -s -A QueueProcessor --max-time 110 'http://127.0.0.1:3331/monitor.php?auth_token=cloudflare_queue_processor_2024' >/dev/null 2>&1 || true; sleep 120; done"
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- seo :3332 — Next.js (фронт). Обращается к бэкенду на localhost:5001 из браузера ---
[program:seo]
directory=/srv/panel/apps/seo
command=/bin/bash -lc "npm run start -- -p 3332"
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3332"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- seo-flask :5001 — бэкенд seo (GSC/OpenAI). Публикуется наружу: фронт зовёт localhost:5001 ---
[program:seo-flask]
directory=/srv/panel/apps/seo
command=/venv/seo/bin/python backend_api.py
autostart=true
autorestart=true
environment=SEO_DATA_DIR="/data/seo"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- ank :3333 — FastAPI + uvicorn (app.main:app), БД в /data/ank ---
[program:ank]
directory=/srv/panel/apps/ank
command=/venv/ank/bin/uvicorn app.main:app --host 0.0.0.0 --port 3333
autostart=true
autorestart=true
environment=DATA_DIR="/data/ank"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- img :3334 — Next + Prisma. На старте: применить схему к volume-БД, засеять, запустить ---
[program:img]
directory=/srv/panel/apps/img
command=/bin/bash -lc "npx prisma db push --skip-generate && npx tsx prisma/seed.ts && npm run start -- -p 3334"
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3334",DATABASE_URL="file:/data/img/imagegen.db",DATA_DIR="/data/img"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- arc :3335 — webarhive (FastAPI factory). Миграции alembic → затем uvicorn ---
[program:arc]
directory=/srv/panel/apps/arc
command=/bin/bash -c "/venv/arc/bin/alembic upgrade head && exec /venv/arc/bin/uvicorn webarhive.web:create_app --factory --host 0.0.0.0 --port 3335 --proxy-headers --forwarded-allow-ips=*"
autostart=true
autorestart=true
environment=DATABASE_URL="sqlite+aiosqlite:////data/arc/webarhive.db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- skins :3336 — Flask + gunicorn (app:app). Сайт на /, панель на /admin ---
[program:skins]
directory=/srv/panel/apps/skins
command=/venv/skins/bin/gunicorn -w 1 -b 0.0.0.0:3336 --timeout 300 app:app
autostart=true
autorestart=true
environment=SITE_ROOT="/srv/panel/apps/skins"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- content :3337 — Content Gap Analyzer (Node/Express). Ключи наследуются из
#     окружения контейнера (.env: OPENAI_API_KEY общий; GOOGLE_NL_API_KEY — свой).
#     Пустой ключ → соответствующий сервис в mock-режиме (приложение всё равно поднимется). ---
[program:content]
directory=/srv/panel/apps/content
command=node src/server.js
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3337"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- mail :3338 — Cloudflare Email Catcher (Node/Express). Читает письма из D1
#     по REST API. Без кредов Cloudflare (CF_ACCOUNT_ID/CF_D1_DATABASE_ID/CF_API_TOKEN
#     в .env или через UI) поднимается в mock-режиме с демо-письмами. Креды, введённые
#     в UI, сохраняются в томе /data/mail (MAIL_DATA_DIR). ---
[program:mail]
directory=/srv/panel/apps/mail
command=node src/server.js
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3338",MAIL_DATA_DIR="/data/mail"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
